# build/conf.d/firefox.conf
[program:firefox]
directory = /opt/momentum
command = firefox
autorestart = false
autostart = false
user = mms
environment = HOME="/opt/momentum",USER="mms"

