# build/conf.d/fluxbox.conf
[program:fluxbox]
directory = /root
command = fluxbox
autorestart = true
autostart = false
user = root
environment = HOME="/root",USER="root"
redirect_stderr = true
