# build/conf.d/lxpanel.conf
[program:lxpanel]
directory = /root
command = /usr/bin/lxpanel --profile LXDE
user = root
environment = HOME="/root",USER="root"
autostart = true
autorestart = true
redirect_stderr = true
