# build/conf.d/momentum.conf
[program:momentum]
directory = /opt/Momentum
command = /opt/Momentum/momentum-prod --no-sandbox
autostart = false
autorestart = true
user = mms
environment = HOME="/opt/momentum",USER="mms"
redirect_stderr = true
