# build/conf.d/openbox.conf
[program:openbox]
directory = /root
command = /usr/bin/openbox
user = root
environment = HOME="/root",USER="root"
autostart = true
autorestart = true
