# build/app/supervisor/conf.d/firefox.conf
[program:firefox]
priority = 30
directory = /opt/momentum
command = /usr/lib/firefox-esr/firefox-esr --profile /opt/momentum/.mozilla
autorestart = false
autostart = false
user = mms
environment = HOME="/opt/momentum",USER="mms"
stderr_events_enabled = true
stdout_events_enabled = true
redirect_stderr = true
