# build/app/supervisor/conf.d/firefox.conf
[program:firefox]
priority = 30
directory = /opt/momentum
command = /usr/lib/firefox-esr/firefox-esr --profile /opt/momentum/.mozilla
autorestart = false
autostart = false
user = mms
environment = HOME="/opt/momentum",USER="mms"
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
