# build/app/supervisor/conf.d/lxpanel.conf
[program:lxpanel]
priority = 15
directory = /root
command = /usr/bin/lxpanel --profile LXDE
user = root
environment = HOME="/root",USER="root"
autorestart = true
stderr_events_enabled = true
stdout_events_enabled = true
redirect_stderr = true
