# build/app/supervisor/conf.d/momentum.conf
[program:momentum]
priority = 30
directory = /opt/Momentum
command = /opt/Momentum/momentum-prod --no-sandbox
autostart = false
autorestart = false
user = mms
environment = HOME="/opt/momentum",USER="mms"
stderr_events_enabled = true
stdout_events_enabled = true
redirect_stderr = true
