# build/app/supervisor/conf.d/nginx.conf
[program:nginx]
priority = 10
command = nginx -c /etc/nginx/nginx.conf -g 'daemon off;'
autorestart = true
directory = /root
stderr_events_enabled = true
stdout_events_enabled = true
redirect_stderr = true
