# build/app/supervisor/conf.d/openbox.conf
[program:openbox]
priority = 15
directory = /root
command = /usr/bin/openbox
user = root
environment = HOME="/root",USER="root"
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
