# build/app/supervisor/conf.d/openbox.conf
[program:openbox]
priority = 15
directory = /root
command = /usr/bin/openbox
user = root
environment = HOME="/root",USER="root"
autostart = true
autorestart = true
stderr_events_enabled = true
stdout_events_enabled = true
redirect_stderr = true
