# build/app/supervisor/conf.d/pcmanfm.conf
[program:pcmanfm]
priority = 15
directory = /root
command = /usr/bin/pcmanfm --desktop --profile LXDE
user = root
stopwaitsecs = 3
user = root
environment = HOME="/root",USER="root"
# autostart = true
autorestart = true
