# build/app/supervisor/conf.d/sshd.conf
[program:sshd]
priority = 30
command = /usr/sbin/sshd -D
autorestart = true
directory = /root
stderr_events_enabled = true
stdout_events_enabled = true
redirect_stderr = true
