# build/app/supervisor/conf.d/vpn-unlimited.conf
[program:vpn-unlimited]
priority = 30
directory = /opt/momentum
command = vpn-unlimited --no-sandbox
autostart = false
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:vpn-unlimited-daemon]
priority = 20
directory = /opt/momentum
command = vpn-unlimited-daemon
autostart = false
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
