# build/app/supervisor/conf.d/websockify.conf
[program:websockify]
priority = 20
command = websockify --web /usr/share/novnc 8080 localhost:5900
autorestart = true
stderr_events_enabled = true
stdout_events_enabled = true
redirect_stderr = true
