# build/app/supervisor/conf.d/x11vnc.conf
[program:x11vnc]
priority = 20
autorestart = true
stderr_events_enabled = true
stdout_events_enabled = true
redirect_stderr = true
command = bash -c "sleep 5 && x11vnc -forever -shared -quiet -rfbauth /.password2"
