# build/app/supervisor/conf.d/x11vnc.conf
[program:x11vnc]
priority = 20
command = x11vnc -forever -shared -quiet  -rfbauth /.password2 -rfbauth /.password2 -rfbauth /.password2 -rfbauth /.password2 -rfbauth /.password2 -rfbauth /.password2 -rfbauth /.password2 -rfbauth /.password2 -rfbauth /.password2 -rfbauth /.password2 -rfbauth /.password2 -rfbauth /.password2 -rfbauth /.password2
autorestart = true
# autostart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
