# build/app/supervisor/conf.d/xterm.conf
[program:xterm]
priority = 20
command = xterm
autorestart = true
autostart = false
user = mms
environment = HOME="/opt/momentum",USER="mms"
stderr_events_enabled = true
stdout_events_enabled = true
redirect_stderr = true
