# build/app/supervisor/conf.d/xterm.conf
[program:xterm]
priority = 20
command = xterm
autorestart = true
autostart = false
user = mms
environment = HOME="/opt/momentum",USER="mms"
